# tweaks/general/profile.d/general.sh
export EDITOR=/usr/bin/micro
export BROWSER=firefox
export PATH=~/.local/bin:$PATH
export TMPDIR='/var/tmp'
