# tweaks/general/profile.d/obs_vkcapture_tweaks.sh
if [ "$XDG_SESSION_TYPE" = "x11" ] && [ $(pacman -Qs obs-vkcapture) != "" ]; then
	export OBS_USE_EGL=1
fi
