# tweaks/profile.d/baloo_disable.sh
if [ $(command -v balooctl) ]; then
	balooctl suspend
	balooctl disable
	balooctl purge
fi
