# tweaks/profile.d/general.sh
export EDITOR=/usr/bin/micro
export BROWSER=firefox
